# service_wallpaper/lightdm_wallpaper.timer
[Unit]
Description=Lightdm wallpaper Timer

[Timer]
OnUnitActiveSec=1day
Unit=lightdm_customize.service

[Install]
WantedBy=timers.target
